# roles/ocserv/templates/ocserv.rc.j2
#!/bin/ksh
# $OpenBSD: ocserv
daemon="/usr/local/sbin/ocserv"
daemon_flags="-c /var/reactance/ocserv/ocserv.conf"

. /etc/rc.d/rc.subr

pexp="ocserv: ocserv-main"

rc_pre() {
	/usr/bin/install -d -o _vpn /var/reactance/ocserv/run/
}

rc_cmd $1
